# scripts/service/caddy/run
#!/bin/sh

set -eu

cd ~/
exec 2>&1
exec caddy -log=stdout -agree=true -conf "${HOME}/.config/caddy/Caddyfile"
