# scripts/service/musikcubed/finish
#!/bin/sh

set -eu

if [ -e /tmp/musikcubed.lock ]; then
	musikcubed --stop || true
	rm -f /tmp/musikcubed.lock
fi
